# wanbot.service
[Unit]
Description=WAN Bot - Ultimate Discord Bot with Web Dashboard
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/wanbot
Environment="PATH=/path/to/wanbot/venv/bin"
ExecStart=/path/to/wanbot/venv/bin/python3 bot.py
Restart=always
RestartSec=10
StandardOutput=append:/path/to/wanbot/bot.log
StandardError=append:/path/to/wanbot/bot_error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/wanbot

[Install]
WantedBy=multi-user.target
